# autostart/rfs.desktop
[Desktop Entry]
Type=Application
Exec=lxterminal --command="/home/pi/dist/main/./main"
